# tests/yara/persistence_sample.service
[Service]
ExecStart=/tmp/.worker
Restart=always

[Install]
WantedBy=multi-user.target
systemctl enable worker.service
