# packaging/systemd/otx-network-monitor.service
[Unit]
Description=OTX-Sec Network Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=OTX_SEC_BASE_DIR=/opt/otx-sec
WorkingDirectory=/opt/otx-sec
ExecStart=/opt/otx-sec/venv/bin/python /opt/otx-sec/monitors/network_monitor.py
Restart=always
RestartSec=10
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
